# ops/com.antonnaumov.codex.telegram-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.antonnaumov.codex.telegram-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>TOKEN="$(/usr/bin/security find-generic-password -s 'codex-telegram-bridge-bot-token' -w 2&gt;/dev/null)" || { echo "Missing Keychain item: codex-telegram-bridge-bot-token" 1&gt;&amp;2; exit 78; }; export CODEX_TELEGRAM_BOT_TOKEN="$TOKEN"; exec /opt/homebrew/bin/node /Users/antonnaumov/code/codex-telegram-frontend/bridge.mjs --config /Users/antonnaumov/code/codex-telegram-frontend/config.local.json</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/antonnaumov</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/antonnaumov/code/codex-telegram-frontend</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/antonnaumov/code/codex-telegram-frontend/logs/bridge.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/antonnaumov/code/codex-telegram-frontend/logs/bridge.stderr.log</string>
</dict>
</plist>
